# homed-gateway-pico/v22.03.5/target/linux/ramips/mt76x8/base-files/etc/uci-defaults/10_set-hostname-and-ssid
. /lib/functions.sh

hostname="HOMEd-GW-Pico"

uci set system.@system[0].hostname=${hostname}
uci commit system

uci set wireless.@wifi-iface[0].ssid=${hostname}
uci set wireless.radio0.disabled='0'
uci commit wireless

wifi
